# snort-machine/block-api.service
[Unit]
Description=Block Management API for SOC Dashboard
After=network.target

[Service]
Type=simple
ExecStart=/usr/bin/python3 /home/ubuntu/snort3/lua/block-api.py
Restart=always
RestartSec=5
User=root

[Install]
WantedBy=multi-user.target
